# config/borg-systemd.service.example
[Unit]
Description=Service for borg-systemd

[Service]
Type=simple
ExecStart=/bin/bash -c 'source /path/to/venv/bin/activate; borg_systemd /path/to/config/config.txt'
