# config/borg-systemd.timer
[Unit]
Description=Run borg-systemd.service every 2 hours

[Timer]
OnStartupSec=2min
OnUnitInactiveSec=2h

[Install]
WantedBy=timers.target
